# roles/router/tasks/conf/iptables-rules.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/iptables-rules.sh

[Install]
WantedBy=default.target
